# 利用cephfs提供k8s的数据存储能力/cephfs_storage-main/cephfs_storage-main/storage/k8s/cephfs/cephfs-provisioner.service
[Unit]
Description=cephfs-provisioner
Documentation=cephfs-provisioner
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cephfs-provisioner -id=cephfs-provisioner-1
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
TimeoutStopSec=20s
RestartSec=10s
StartLimitInterval=30min
StartLimitBurst=3


[Install]
WantedBy=multi-user.target
